# automation/com.wealth.digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Weekly digest, Saturdays at 08:00. Install: see SETUP.md. Paths are already filled in. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>            <string>com.wealth.digest</string>
  <key>ProgramArguments</key>
  <array>
    <string>/PATH/TO/wealth-dashboard/automation/.venv/bin/python3</string>
    <string>/PATH/TO/wealth-dashboard/automation/digest.py</string>
  </array>
  <key>WorkingDirectory</key> <string>/PATH/TO/wealth-dashboard/automation</string>
  <key>StartCalendarInterval</key>
  <dict><key>Weekday</key><integer>6</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
  <key>RunAtLoad</key>        <false/>
  <key>StandardOutPath</key>  <string>/PATH/TO/wealth-dashboard/automation/digest.log</string>
  <key>StandardErrorPath</key><string>/PATH/TO/wealth-dashboard/automation/digest.log</string>
</dict>
</plist>
